# workspace7/crontab/multisites
* * * * * laradock php /var/www/multisites/artisan schedule:run >> /dev/null 2>&1

